# dev/.scripts/project-zfs_utils/rc.local
# "rpool/ROOT" may need to be manually modified
zfs list -H -o name -r rpool/ROOT |
        xargs -i@ zfs set mountpoint=/@ @ && true
zfs mount -a && true
